# files/etc/hotplug.d/iface/50-mesh-tx-offload
#!/bin/ash
ethtool --offload phy2-mesh0 tx-checksumming off